# ops/launchd/st.urm.brain-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>st.urm.brain-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/steveu/code/brain-sync/bin/brain-sync</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>WatchPaths</key>
    <array>
        <string>/Users/steveu/brain/vault</string>
        <string>/Users/steveu/Library/Mobile Documents/iCloud~md~obsidian/Documents/Brain</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/steveu</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/steveu/Library/Logs/brain-sync.launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/steveu/Library/Logs/brain-sync.launchd.log</string>
</dict>
</plist>
